# prebuilt/fs.install/etc/inittab
version https://git-lfs.github.com/spec/v1
oid sha256:3e66a16ce2fcdca776587b28976f334efbc0d72ec0217b31ae3d325564d876c9
size 425
